# ops/ansible/roles/factory_deploy/templates/yamiochi-factory.service.j2
[Unit]
Description=Yamiochi Fabro server
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
WorkingDirectory={{ factory_root }}
ExecStartPre=-/usr/bin/docker rm -f yamiochi-fabro
ExecStartPre=/usr/bin/docker pull {{ factory_fabro_image }}
ExecStart=/usr/bin/docker run --name yamiochi-fabro --rm --user root --network host --env-file {{ factory_root }}/.env.runtime --env-file {{ factory_root }}/.env.server -v {{ factory_state_root }}:{{ factory_state_root }} -v {{ factory_state_root }}/fabro-storage:/storage -v {{ factory_root }}/config:/config:ro -v /var/run/docker.sock:/var/run/docker.sock --entrypoint fabro {{ factory_fabro_image }} server start --foreground --bind 0.0.0.0:{{ factory_fabro.listen_port }} --config /config/settings.toml
ExecStop=/usr/bin/docker stop yamiochi-fabro
Restart=always
RestartSec=5
TimeoutStartSec=300
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
